# deploy/systemd/aman-queue.service
# Laravel queue worker (QUEUE_CONNECTION=database)
# Install:  cp deploy/systemd/aman-queue.service /etc/systemd/system/
#           systemctl enable --now aman-queue
[Unit]
Description=Aman Laravel queue worker
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/var/www/aman/backend
ExecStart=/usr/bin/php /var/www/aman/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600

[Install]
WantedBy=multi-user.target
